# docker/deploy/supervisord-backend.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:guacd]
command=/usr/local/sbin/guacd -f -b 127.0.0.1 -l 4822
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:backend]
; Run as root so it can read guacd recording files (guacd creates 0640 root-owned files)
command=node /app/apps/backend/dist/index.js
directory=/app
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=NODE_ENV="production"
